# package/DEBIAN/postinst
#!/bin/sh
/bin/systemctl enable tailscalesd.service
/bin/systemctl start tailscalesd.service
