# package/DEBIAN/prerm
#!/bin/sh
/bin/systemctl stop tailscalesd.service
/bin/systemctl disable tailscalesd.service
